# deploy/openage.service
[Unit]
Description=OpenAge 生物年龄预测服务
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/openage
Environment="ANTHROPIC_API_KEY=your_key_here"
ExecStart=/opt/openage/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
